# ops/systemd/aideal-free-llm-health-probe.timer
[Unit]
Description=Run AIdeal CPS Free LLM Health Probe daily

[Timer]
OnCalendar=*-*-* 04:20:00
Persistent=true
RandomizedDelaySec=20m
AccuracySec=5m
Unit=aideal-free-llm-health-probe.service

[Install]
WantedBy=timers.target
