# ops/systemd/aideal-catalog-refresh.timer
[Unit]
Description=Run AIdeal CPS Catalog Refresh nightly

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=aideal-catalog-refresh.service

[Install]
WantedBy=timers.target
