# ops/systemd/aideal-free-llm-health-probe.timer
[Unit]
Description=Run AIdeal CPS Free LLM Health Probe every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=aideal-free-llm-health-probe.service

[Install]
WantedBy=timers.target
